# server/karkas-leads.service
[Unit]
Description=Karkas Invest Telegram lead receiver
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=karkas-leads
Group=karkas-leads
WorkingDirectory=/opt/karkas-leads
EnvironmentFile=/etc/karkas-leads.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/karkas-leads/app.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
